# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[program:zep_uwsgi]
command=/home/ubuntu/zenv/bin/uwsgi --ini /home/ubuntu/zep_tunes/zep.ini
autostart=true
autorestart=true
user=ubuntu
directory=/home/ubuntu/zep_tunes/
stopwaitsecs=1
startretries=60

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock